# plugin/app/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:docker-gen]
command=docker-gen -watch -notify "crontab /app/crontab.txt" /app/crontab.tmpl /app/crontab.txt
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:cron]
command=/usr/sbin/crond -f -l 2
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
